# deploy/bridge-crm.service
[Unit]
Description=Bridge CRM gunicorn service
After=network.target postgresql.service

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/inventory-chatbot
EnvironmentFile=/home/ubuntu/inventory-chatbot/bridge_crm/.env
ExecStart=/home/ubuntu/crm-env/bin/gunicorn \
  --workers 3 \
  --bind 127.0.0.1:5001 \
  --access-logfile - \
  --error-logfile - \
  bridge_crm.wsgi:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
